# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
loglevel=info

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err.log
stdout_logfile=/var/log/supervisor/nginx.out.log
priority=10

[program:flask]
command=/app/venv/bin/python /app/backend/run.py
directory=/app/backend
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/flask.err.log
stdout_logfile=/var/log/supervisor/flask.out.log
environment=FLASK_ENV="production",DATABASE_PATH="/app/storage/database/funhub.db",THUMBNAIL_FOLDER="/app/storage/thumbnails",CORS_ORIGINS="http://localhost,http://127.0.0.1"
priority=20
